# etc/systemd/system/benchmark-setup.service
[Unit]
Description=Configure services for benchmarking

[Service]
Type=oneshot
ExecStart=/opt/benchmark/setup.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
